# deploy/systemd/futu-opend.service
[Unit]
Description=Futu OpenD
After=network.target

[Service]
Type=simple
User=futu
Group=futu
WorkingDirectory=/opt/futu-opend
ExecStart=/opt/futu-opend/FutuOpenD -c /opt/futu-opend/OpenD.xml -no_monitor=1 -auto_hold_quote_right=1
Restart=always
RestartSec=3
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
