# deploy/systemd/hk-tick-collector.service
[Unit]
Description=HK Tick Collector (Futu OpenD -> SQLite WAL)
After=network-online.target futu-opend.service
Wants=network-online.target
Requires=futu-opend.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=hkcollector
Group=hkcollector
WorkingDirectory=/opt/futu_tick_downloader
EnvironmentFile=-/etc/hk-tick-collector.env
ExecStart=/opt/futu_tick_downloader/.venv/bin/python -m hk_tick_collector.main
Restart=always
RestartSec=3
TimeoutStartSec=30
TimeoutStopSec=120
KillSignal=SIGINT
LimitNOFILE=65535
Nice=5
UMask=0027

# Security hardening (保守配置，避免鎖死維運)
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
ProtectHome=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectKernelLogs=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
ReadWritePaths=/data/sqlite/HK /opt/futu_tick_downloader

[Install]
WantedBy=multi-user.target
