# deploy/systemd/hk-tick-collector.service
[Unit]
Description=HK Tick Collector (Futu OpenD to SQLite)
After=network-online.target futu-opend.service
Wants=network-online.target
Requires=futu-opend.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=hkcollector
Group=hkcollector
WorkingDirectory=/opt/futu_tick_downloader
EnvironmentFile=/opt/futu_tick_downloader/.env
ExecStart=/opt/futu_tick_downloader/.venv/bin/python -m hk_tick_collector.main
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=180
LimitNOFILE=65535
UMask=0027

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
ReadWritePaths=/data/sqlite/HK /opt/futu_tick_downloader

[Install]
WantedBy=multi-user.target
