# ops/hk-tick-collector.service
[Unit]
Description=HK Tick Collector
After=network.target futu-opend.service
Requires=futu-opend.service
# If configuring StartLimitIntervalSec, keep it under [Unit], not [Service].

[Service]
Type=simple
User=hkcollector
Group=hkcollector
WorkingDirectory=/opt/futu_tick_downloader
# If you prefer repo-local env, set EnvironmentFile=/opt/futu_tick_downloader/.env
EnvironmentFile=/etc/hk-tick-collector.env
ExecStart=/opt/futu_tick_downloader/.venv/bin/python -m hk_tick_collector.main
Restart=always
RestartSec=5
TimeoutStopSec=120
KillSignal=SIGINT
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
